# cups/rootfs/etc/services.d/cups/run
#!/usr/bin/with-contenv bashio

#Setup authentication
CUPSADMIN=$(bashio::config 'user')
CUPSPASSWORD=$(bashio::config 'passwd')



if [ $(grep -ci $CUPSADMIN /etc/shadow) -eq 0 ]; then
    adduser -S -G lpadmin -H $CUPSADMIN 
fi
echo $CUPSADMIN:$CUPSPASSWORD | chpasswd


# copy config
cp /data/config/printers.conf /etc/cups/printers.conf
/usr/sbin/cupsd -t
/usr/sbin/cupsd 
/opt/printer-update.sh
